# scripts/tuition-payment-monitor.service
[Unit]
Description=Tuition Payment System Monitoring Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/bin/bash /app/scripts/monitor.sh

[Install]
WantedBy=multi-user.target